# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:zookeeper]
command=/usr/bin/start-zookeeper.sh

[program:kafka]
command=/usr/bin/start-kafka.sh
